# services/hddtemp.service
[Unit]
Description=HDD temperature monitor daemon

[Service]
ExecStart=/usr/bin/hddtemp -dF /dev/sdb -p 7635

[Install]
WantedBy=multi-user.target
