# services/hddtemp.service
[Unit]
Description=HDD temperature monitor daemon

[Service]
ExecStart=/usr/bin/hddtemp -dF /dev/sda -p 7634

[Install]
WantedBy=multi-user.target
